# deploy/prod/livekit-server.service
[Unit]
Description=LiveKit SFU
After=network.target

[Service]
Type=simple
User=openclaw
ExecStart=/usr/local/bin/livekit-server --config /etc/livekit/livekit.yaml
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

LimitNOFILE=500000

[Install]
WantedBy=multi-user.target
